# systemd/guardian.service
[Unit]
Description=Guardian Embedded Health Monitor
After=network.target

[Service]
ExecStart=/usr/local/bin/guardian
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
